# roles/serial/templates/colorado.service
[Unit]
Description=colorado
After=network.target

[Service]
ExecStart=/usr/bin/colorado
WorkingDirectory=/opt/colorado
StandardOutput=inherit
StandardError=inherit
Restart=always
User=root

[Install]
WantedBy=multi-user.target